# tklpatch/musicbrainz/overlay/etc/cron.daily/musicbrainz_loadRepPackets
#!/bin/sh

DIRECTORY=`cat /etc/MUSICBRAINZ`
cd $DIRECTORY
cd admin/replication

./LoadReplicationChanges > /var/log/musicbrainz_loadRepChanges.log

